# swarm/customerAPI-swarm-cent.service
Description = Test Service by techken.in
After = network.target
 
[Service]
ExecStart = /root/scripts/myservice.sh
 
[Install]
WantedBy = multi-user.target